# he-ddns-updater_1.0-1/etc/systemd/system/he-ddns-updater.timer
[Unit]
Description=Run he ddns updater hourly

[Timer]
OnCalendar=hourly
RandomizedDelaySec=1min
Persistent=true

[Install]
WantedBy=timers.target
