# he-ddns-updater_1.0-1/etc/systemd/system/he-ddns-updater.service
[Unit]
Description=HE DDNS Updater
[Service]
Type=oneshot
ExecStart=/usr/bin/update-he-ddns
[Install]
WantedBy=multi-user.target
