# he-ddns-updater_1.0-1/DEBIAN/prerm
#!/bin/sh

systemctl stop he-ddns-updater.timer || true
systemctl disable he-ddns-updater.timer || true
